# assets/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
priority = 1
command = /usr/sbin/php-fpm7.1 -F
autostart = true
user = root
redirect_stderr=true
stdout_logfile = /var/log/supervisor/php-fpm.log
loglevel = warn
